# server/task.cron
* * * * * /usr/bin/python /app/manage.py check_paying
* * * * * /usr/bin/python /app/manage.py check_expired
